# com.saul.milltender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.saul.milltender</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/saul/src/treadmill/.venv/bin/python</string>
        <string>-u</string>
        <string>/Users/saul/src/treadmill/milltender.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/saul/src/treadmill</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/saul/src/treadmill/milltender.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/saul/src/treadmill/milltender.log</string>
</dict>
</plist>
